# btrfs-dwarfs/configs/bdfs_daemon.service.in
# SPDX-License-Identifier: GPL-2.0-or-later
# systemd service unit for the BTRFS+DwarFS daemon
# Generated by cmake — do not edit directly; edit configs/bdfs_daemon.service.in

[Unit]
Description=BTRFS+DwarFS hybrid filesystem daemon
Documentation=man:bdfs_daemon(8)
After=local-fs.target
Requires=local-fs.target

[Service]
Type=forking
PIDFile=/run/bdfs/daemon.pid
ExecStartPre=/bin/mkdir -p /run/bdfs /var/lib/bdfs
ExecStart=@BDFS_DAEMON_BIN@ \
    -c /dev/bdfs_ctl \
    -s /run/bdfs/daemon.sock \
    -d /var/lib/bdfs \
    -j 4
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5s

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/run/bdfs /var/lib/bdfs
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
